# deploy/host/openclaw-mail-bff.service.example
[Unit]
Description=Mail Agent BFF
After=network.target

[Service]
Type=simple
User=mail-agent
Group=mail-agent
WorkingDirectory=/opt/mail-agent-bff
Environment=NODE_ENV=production
EnvironmentFile=/etc/mail-agent-bff/bff.env
ExecStart=/usr/bin/env node /opt/mail-agent-bff/apps/bff/dist/server.js
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/mail-agent-bff/apps/bff/data

[Install]
WantedBy=multi-user.target
